# deploy/poster.timer
[Unit]
Description=KMotors — тик автопостинга авто (каждые 2 часа)

[Timer]
# Тикаем каждые 2 часа. Окно публикации (9–21 KST) и дневной лимит (5/день)
# проверяются внутри роута /api/poster/run, поэтому вне окна тик просто ничего
# не постит. Джиттер до 20 мин — чтобы посты не выходили ровно по часам.
OnCalendar=*-*-* 00/2:00:00
Persistent=true
RandomizedDelaySec=1200

[Install]
WantedBy=timers.target
